# files/system/etc/systemd/system/zram-recompression.timer
[Unit]
Description=Control timing to recompress zram devices

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
RandomizedDelaySec=10s
Unit=zram-recompression.service

[Install]
WantedBy=multi-user.target
