# files/system/etc/systemd/system/zram-recompression.service
[Unit]
Description=Recompress idle+huge, then trigger zram compaction
Documentation="https://docs.kernel.org/admin-guide/blockdev/zram.html#recompression"
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/recompress-zram
