# files/system/etc/systemd/system/zram-recompression.service
[Unit]
Description=Mark all pages as idle, recompress idle+huge pages, then trigger compaction.
Documentation="https://docs.kernel.org/admin-guide/blockdev/zram.html#recompression"
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/recompress-zram
